# docker/etc/cron.d/changepwds
* * * * * root python /root/scripts/password_reset.py
