# docker/etc/cron.d/myjob
* * * * * root echo 'this is some text' >> /tmp/myjob.log
