# docker/etc/cron.d/print
* * * * * root python /scripts/print_stuff.py >> /tmp/print.log
