# docker/var/spool/cron/crontabs/root
* * * * * python3 /root/scripts/password_reset.py
